# infra/roles/docker/templates/docker.service.j2
# {{ ansible_managed }}

[Unit]
Description=Docker configuration for Kubernetes
Documentation=https://docs.docker.com/install/linux/linux-postinstall/
After=network-online.target docker.socket firewalld.service
Wants=network-online.target
Requires=docker.socket
 
[Service]
Type=notify
# Docs
# https://kubernetes.io/docs/setup/scratch/#docker
# https://docs.docker.com/install/linux/linux-postinstall/#configuring-remote-access-with-systemd-unit-file
# https://coreos.com/flannel/docs/latest/flannel-config.html
# https://icicimov.github.io/blog/kubernetes/Kubernetes-cluster-step-by-step-Part4/
# FlannelD subnet setup
# the following file is generated by Flannel incl the Environment variables
EnvironmentFile=/var/lib/flanneld/subnet.env
ExecStart=/usr/bin/dockerd \
    --bip=${FLANNEL_SUBNET} \
    --mtu=${FLANNEL_MTU} \
    --iptables=false \
    --ip-masq=false \
    --ip-forward=true \
    -H fd://

ExecReload=/bin/kill -s HUP $MAINPID
LimitNOFILE=1048576
# Having non-zero Limit*s causes performance problems due to accounting overhead
# in the kernel. We recommend using cgroups to do container-local accounting.
LimitNPROC=infinity
LimitCORE=infinity
# Uncomment TasksMax if your systemd version supports it.
# Only systemd 226 and above support this version.
TasksMax=infinity
TimeoutStartSec=0
# set delegate yes so that systemd does not reset the cgroups of docker containers
Delegate=yes
# kill only the docker process, not all processes in the cgroup
KillMode=process
# restart the docker process if it exits prematurely
Restart=on-failure
StartLimitBurst=3
StartLimitInterval=60s

[Install]
WantedBy=multi-user.target
